# files/etc/uci-defaults/01-system
#!/bin/sh

uci set system.@system[0].hostname='Office-Gateway'

uci commit system

exit 0
